# scripts/com.shenjiying88.nightly-foundation11.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
  <key>Label</key>
  <string>com.shenjiying88.nightly-foundation11</string>
  <key>ProgramArguments</key>
  <array>
    <string>/bin/bash</string>
    <string>/Users/yaoyunzhong/Desktop/shenjiying/shenjiying88/scripts/nightly-jobs.sh</string>
  </array>
  <key>RunAtLoad</key>
  <false/>
  <key>StartCalendarInterval</key>
  <dict>
    <key>Hour</key>
    <integer>0</integer>
    <key>Minute</key>
    <integer>0</integer>
  </dict>
  <key>WorkingDirectory</key>
  <string>/Users/yaoyunzhong/Desktop/shenjiying/shenjiying88</string>
  <key>StandardOutPath</key>
  <string>/Users/yaoyunzhong/Desktop/shenjiying/shenjiying88/docs/monitoring/nightly/launchd.log</string>
  <key>StandardErrorPath</key>
  <string>/Users/yaoyunzhong/Desktop/shenjiying/shenjiying88/docs/monitoring/nightly/launchd-errors.log</string>
  <key>LowPriorityIO</key>
  <true/>
  <key>ProcessType</key>
  <string>Background</string>
</dict>
</plist>
